# nagios/asm_nagios.cron
*/5 * * * * root  [ -x /opt/asm-deployer/nagios/poll_idrac8_sensors.rb ] && /usr/lib64/nagios/plugins/load_short_circuit.sh -s -l 4 -e 1 && /opt/asm-deployer/nagios/poll_idrac8_sensors.rb >> /opt/Dell/ASM/logs/poll_idrac8_sensors.log 2>&1
*/2 * * * * root  [ -x /opt/asm-deployer/nagios/nagios-export.rb ] && /usr/lib64/nagios/plugins/load_short_circuit.sh -s -e 1 && /opt/asm-deployer/nagios/nagios-export.rb > /opt/Dell/ASM/logs/nagios-export.log 2>&1
* * * * * root  [ -x /usr/lib64/nagios/plugins/process_nagios_status.rb ] && /usr/lib64/nagios/plugins/load_short_circuit.sh -s -e 1 && /usr/lib64/nagios/plugins/process_nagios_status.rb -d > /opt/Dell/ASM/logs/process_nagios_status.log 2>&1
